# roles/ledger/templates/systemd/ledger-networth.timer.j2
{{ ansible_managed | comment }}

[Unit]
Description=Add recurring networth calculation to Ledger

[Timer]
OnCalendar=*-*-* 23:50:00

[Install]
WantedBy=timers.target
