# roles/ledger/templates/systemd/ledger-sync.timer.j2
{{ ansible_managed | comment }}

[Unit]
Description=Sync local and remote changes with Firefly everyday

[Timer]
OnCalendar=*-*-* 03:00:00

[Install]
WantedBy=timers.target
